# config/deploy/recipes/templates/init.d/default.erb
#! /bin/sh
### BEGIN INIT INFO
# Provides:		    <%= name %>
# Required-Start:	$syslog $remote_fs
# Required-Stop:	$syslog $remote_fs
# Should-Start:		$local_fs
# Should-Stop:		$local_fs
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
### END INIT INFO

PATH=<%= path %>
DAEMON=<%= daemon %>
DAEMON_ARGS=<%= daemon_args if defined?(daemon_args) %>
NAME=<%= name %>
DESC=<%= name %>

RUNDIR=<%= defined?(pidfile) && pidfile ? File.dirname(pidfile) : "/var/run/#{name}" %>
PIDFILE=<%= defined?(pidfile) && pidfile ? pidfile : "$RUNDIR/#{name}.pid" %>

test -x $DAEMON || exit 0

set -e

case "$1" in
  start)
  echo -n "Starting $DESC: "
  mkdir -p $RUNDIR
  touch $PIDFILE
  chmod 755 $RUNDIR
  if start-stop-daemon --start --quiet --umask 007 --pidfile $PIDFILE --exec $DAEMON -- $DAEMON_ARGS
  then
    echo "success"
  else
    echo "failed"
  fi
  ;;
  stop)
  echo -n "Stopping $DESC: "
  if start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE
  then
    echo "success"
  else
    echo "failed"
  fi
  rm -f $PIDFILE
  ;;

  restart|force-reload)
  ${0} stop
  ${0} start
  ;;

  status)
  echo -n "$DESC is "
  if start-stop-daemon --stop --quiet --signal 0 --name ${NAME} --pidfile ${PIDFILE}
  then
    echo "running"
  else
    echo "not running"
    exit 1
  fi
  ;;

  *)
  echo "Usage: /etc/init.d/$NAME {start|stop|restart|force-reload}" >&2
  exit 1
  ;;
esac

exit 0